# systemd/hermes-canvas-server-dev.service
[Unit]
Description=Hermes Canvas gateway (development)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/hermes-canvas/hermes-canvas.env
Environment=PATH=__HERMES_CANVAS_NODE_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=__HERMES_CANVAS_ROOT__
ExecStart=__HERMES_CANVAS_PACKAGE_MANAGER__ run server
Restart=on-failure
RestartSec=5
SyslogIdentifier=hermes-canvas-server-dev

[Install]
WantedBy=multi-user.target
