# systemd/hermes-canvas-app.service
[Unit]
Description=Hermes Canvas frontend
After=network.target hermes-canvas-server.service
Wants=hermes-canvas-server.service

[Service]
Type=simple
EnvironmentFile=/etc/hermes-canvas/hermes-canvas.env
Environment=PATH=__HERMES_CANVAS_NODE_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=__HERMES_CANVAS_ROOT__
ExecStart=__HERMES_CANVAS_PACKAGE_MANAGER__ run serve:app
Restart=on-failure
RestartSec=5
SyslogIdentifier=hermes-canvas-app

[Install]
WantedBy=multi-user.target
